# docker/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
user = root
autostart = true

[program:php5-fpm]
command = /usr/sbin/php5-fpm -R
user = root
autostart = true
